# launchagents/com.obsidian.vault-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian.vault-health</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>__HOME__/claude-vault/vault-health.py</string>
    </array>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__HOME__/claude-vault/vault-health-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/claude-vault/vault-health-stderr.log</string>
</dict>
</plist>
