# launchagents/com.obsidian.topic-linker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.obsidian.topic-linker</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>__HOME__/claude-vault/topic-linker.py</string>
    </array>
    <key>StartInterval</key>
    <integer>600</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <false/>
    <key>StandardOutPath</key>
    <string>__HOME__/claude-vault/topic-linker-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>__HOME__/claude-vault/topic-linker-stderr.log</string>
</dict>
</plist>
